# infra/supervisor/programs/accumulo-tserver.ini
[program:accumulo-tserver]
command=/bin/sh -lc 'exec accumulo tserver'
autostart=true
autorestart=true
stdout_logfile=/Volumes/Samsung_T5/var/log/hazoo/accumulo-tserver.out
stderr_logfile=/Volumes/Samsung_T5/var/log/hazoo/accumulo-tserver.err
environment=JAVA_HOME="/Library/Java/JavaVirtualMachines/temurin-17.jdk/Contents/Home", ACCUMULO_HOME="/Volumes/Samsung_T5/accumulo", PATH="/Library/Java/JavaVirtualMachines/temurin-17.jdk/Contents/Home/bin:/Volumes/Samsung_T5/accumulo/bin:%(ENV_PATH)s"
