# infra/supervisor/programs/accumulo-monitor.ini
[program:accumulo-monitor]
command=/bin/sh -lc 'exec accumulo monitor'
autostart=true
autorestart=true
stdout_logfile=/Volumes/Samsung_T5/var/log/hazoo/accumulo-monitor.out
stderr_logfile=/Volumes/Samsung_T5/var/log/hazoo/accumulo-monitor.err
environment=JAVA_HOME="/Library/Java/JavaVirtualMachines/temurin-17.jdk/Contents/Home", ACCUMULO_HOME="/Volumes/Samsung_T5/accumulo", PATH="/Library/Java/JavaVirtualMachines/temurin-17.jdk/Contents/Home/bin:/Volumes/Samsung_T5/accumulo/bin:%(ENV_PATH)s"
