# infra/supervisor/programs/haz-microservice.ini
[program:haz-microservice]
directory=/Volumes/Samsung_T5/d4m-microservice
command=/Library/Java/JavaVirtualMachines/temurin-17.jdk/Contents/Home/bin/java -jar /Volumes/Samsung_T5/d4m-microservice/d4m.acc.microservice-0.0.1.jar --server.port=5102 --accumulo.client.props=file:///Volumes/Samsung_T5/accumulo/conf/accumulo-client.properties
autostart=false
autorestart=true
stdout_logfile=/Volumes/Samsung_T5/var/log/hazoo/haz-microservice.out
stderr_logfile=/Volumes/Samsung_T5/var/log/hazoo/haz-microservice.err
